# src/ol_configuration_management/components/concourse/templates/concourse_web.service.j2
[Unit]
Description=Concourse service for managing continuous deployment pipelines
After=network.target

[Service]

EnvironmentFile=/opt/concourse/env

{% if worker %}
KillSignal=SIGUSR2
{% endif %}
